# packaging/debian-template/DEBIAN/prerm
#!/bin/bash
set -e

# Stop the service if it's running
if systemctl is-active --quiet mc-daemon.service; then
    echo "Stopping mc-daemon service..."
    systemctl stop mc-daemon.service
fi

# Disable the service from starting on boot
if systemctl is-enabled --quiet mc-daemon.service; then
    echo "Disabling mc-daemon service..."
    systemctl disable mc-daemon.service
fi

exit 0
