# Doc/Samples/ansible/templates/mc-daemon.service.j2
[Unit]
Description=Meadow Daemon - OTA Update and Cloud Connectivity Service
Documentation=https://github.com/WildernessLabs/Meadow.Daemon
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/meadow.conf

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/mc-daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mc-daemon

# Security hardening (optional - adjust as needed)
# Note: Running as root is required for system operations
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ProtectHome=false

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
